# scripts/launchd/com.amswiki.kakao-cookie-refresh.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  카카오 파트너센터 쿠키 선제 갱신 (6시간마다). Chrome 로컬 쿠키 → .env.local.
  쿠키 만료로 수집 데몬이 멈추는 것을 예방. 값이 바뀐 경우에만 .env.local 을 갱신하며,
  KAKAO_SKIP_DAEMON_RESTART=1 이라 수집 데몬을 끊지 않는다(데몬은 다음 401 때 .env.local 재읽음).
  전제: Chrome 으로 business.kakao.com 에 로그인 유지 + 최초 1회 "Chrome Safe Storage" 키체인 허용.
  설치:
    1) WorkingDirectory / node 경로 / 로그 경로를 본인 머신에 맞게 수정
    2) cp scripts/launchd/com.amswiki.kakao-cookie-refresh.plist ~/Library/LaunchAgents/
    3) launchctl load ~/Library/LaunchAgents/com.amswiki.kakao-cookie-refresh.plist
    4) tail -f ~/Library/Logs/ams-wiki/kakao-cookie-refresh.log
  중지:
    launchctl unload ~/Library/LaunchAgents/com.amswiki.kakao-cookie-refresh.plist
-->
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.amswiki.kakao-cookie-refresh</string>

  <key>ProgramArguments</key>
  <array>
    <string>/usr/local/bin/node</string>
    <string>scripts/kakao-partner-refresh-cookie.mjs</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/layr8x/Library/Mobile Documents/com~apple~CloudDocs/🕹️ layr8x/ams-wiki</string>

  <key>RunAtLoad</key>
  <true/>

  <key>StartInterval</key>
  <integer>21600</integer>

  <key>StandardOutPath</key>
  <string>/Users/layr8x/Library/Logs/ams-wiki/kakao-cookie-refresh.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/layr8x/Library/Logs/ams-wiki/kakao-cookie-refresh.err.log</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>NODE_ENV</key>
    <string>production</string>
    <key>KAKAO_SKIP_DAEMON_RESTART</key>
    <string>1</string>
  </dict>
</dict>
</plist>
