# scripts/launchd/com.amswiki.kakao-stream.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  카카오 파트너센터 3채널(_VGAQn,_TkpPG,_xfxilXn) 실시간 수집 데몬을 macOS launchd 로 상시 실행.
  multi-stream supervisor 가 채널마다 자식 프로세스를 띄움. 경로는 이 머신 기준으로 채워져 있음.
  설치:
    1) cp scripts/launchd/com.amswiki.kakao-stream.plist ~/Library/LaunchAgents/
    2) launchctl load ~/Library/LaunchAgents/com.amswiki.kakao-stream.plist
    3) tail -f ~/Library/Logs/ams-wiki/kakao-stream.log
  중지:
    launchctl unload ~/Library/LaunchAgents/com.amswiki.kakao-stream.plist
  쿠키 갱신 후 재시작:
    launchctl kickstart -k gui/$(id -u)/com.amswiki.kakao-stream
-->
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.amswiki.kakao-stream</string>

  <key>ProgramArguments</key>
  <array>
    <string>/usr/local/bin/node</string>
    <string>--env-file=.env.local</string>
    <string>scripts/kakao-partner-multi-stream.mjs</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/layr8x/Library/Mobile Documents/com~apple~CloudDocs/🕹️ layr8x/ams-wiki</string>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <dict>
    <key>SuccessfulExit</key>
    <false/>
    <key>Crashed</key>
    <true/>
  </dict>

  <key>ThrottleInterval</key>
  <integer>30</integer>

  <key>StandardOutPath</key>
  <string>/Users/layr8x/Library/Logs/ams-wiki/kakao-stream.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/layr8x/Library/Logs/ams-wiki/kakao-stream.err.log</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>NODE_ENV</key>
    <string>production</string>
  </dict>
</dict>
</plist>
